# owncloud-server/etc/cron.d/clamAV
52 * * * * /usr/bin/freshclam --quiet #Update clamAV virus db
