# owncloud-server/etc/cron.d/oc-rescanfiles
*/15 * * * * www-data /usr/bin/php /var/www/owncloud/occ files:scan --all

