# owncloud-server/etc/cron.d/oc-rescanfiles
*/15 * * * * www-data php occ files:scan --all
