# postgis/etc/cron.d/backupPGSQL
02 01 * * * /root/backupPGSQL.sh # Performs a daily full backup of the PGSQL DBs, under /backups/
